# files/etc/uci-defaults/99-setup
#!/bin/sh
LOG=/tmp/r4s-bridge.log
exec >"$LOG" 2>&1

echo "=== R4S BRIDGE START $(date) ==="

# =========================
# 1. IPv6 ULA 清理（可选但安全）
# =========================
uci -q delete network.globals.ula_prefix

# =========================
# 2. LAN 桥接（核心）
# =========================
uci set network.lan.device='br-lan'

# 清理旧端口（避免重复）
uci -q del_list network.lan.ports='eth0'
uci -q del_list network.lan.ports='eth1'

# 加入双口
uci add_list network.lan.ports='eth0'
uci add_list network.lan.ports='eth1'

# =========================
# 3. 管理IP
# =========================
uci set network.lan.ipaddr='192.168.1.2'
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.gateway='192.168.1.1'
uci set network.lan.dns='192.168.1.1'

# =========================
# 4. DHCP关闭（交给光猫）
# =========================
uci set dhcp.lan.ignore='1'
uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.dhcpv6='disabled'

# =========================
# 5. 提交配置
# =========================
uci commit network
uci commit dhcp

# =========================
# 6. 关闭IPv6（最稳方式）
# =========================
echo 1 > /proc/sys/net/ipv6/conf/all/disable_ipv6

# =========================
# 7. 重启网络
# =========================
/etc/init.d/network restart

echo "=== DONE $(date) ==="
